# shell/inittab
::sysinit:/usr/bin/busybox --install -s
::respawn:/bin/sh
